# scripts/lib/systemd/system/robolab-renderer-frontend.service
[Unit]
Description=Robolab renderer frontend (stable)

[Service]
Type=exec
Restart=always
User=deploy
ExecStart=node index.js
WorkingDirectory=/home/deploy/web-projects/robolab-renderer/stable/frontendServer

[Install]
WantedBy=multi-user.target
